# confs/supervisor/bcexd.conf
[program:bcexd]
command = bcexchanged -datadir=/home/bcex/.bcexchange
user = bcex
redirect_stderr = true
stdout_logfile = /home/bcex/bcexchange.log

