# confs/supervisor/charts.conf
[program:charts]
directory=/home/bcex/BCExplorer_OS
user=bcex
command=/home/bcex/BCExplorer_OS/charts.bash
redirect_stderr=True
stdout_logfile=/home/bcex/charts.log
