# confs/supervisor/explorer.conf
[program:explorer]
user = bcex
directory = /home/bcex/BCExplorer_OS
command = grunt server
redirect_stderr = true
stdout_logfile = /home/bcex/explorer.log
