# confs/supervisor/parser.conf
[program:parser]
directory=/home/bcex/BCExplorer_OS
user=bcex
command=/home/bcex/BCExplorer_OS/runparser.bash
redirect_stderr=True
stdout_logfile=/home/bcex/parser.log
